# images/stackstorm/supervisor/mistral.conf
[program:mistral]
process_name = %(program_name)s
command = /etc/init.d/%(program_name)s start
numprocs = 1
process_name = %(program_name)s
redirect_stderr = true
