# images/stackstorm/supervisor/st2chatops.conf
[program:st2chatops]
user=st2
group=st2
directory=/opt/stackstorm/chatops
command=/opt/stackstorm/chatops/bin/hubot ''
