# images/stackstorm/supervisor/st2notifier.conf
[program:st2notifier]
user=st2
group=st2
command=/opt/stackstorm/st2/bin/st2notifier --config-file /etc/st2/st2.conf
