# images/stackstorm/supervisor/st2rulesengine.conf
[program:st2rulesengine]
user=st2
group=st2
command=/opt/stackstorm/st2/bin/st2rulesengine --config-file /etc/st2/st2.conf
