# images/stackstorm/supervisor/st2timersengine.conf
[program:st2timersengine]
user=st2
group=st2
command=/opt/stackstorm/st2/bin/st2timersengine --config-file /etc/st2/st2.conf
